# config/skydns.upstart
description "skydns service registry"

start on started
stop on shutdown

exec /usr/bin/skydns